# cron.d/scrape
# m h dom mon dow user  command
15 9	* * *	slgmibu	~/githubworkspace/Lake-Water-Level/scrape.pl
