# cron.d/scrape
# m h dom mon dow user  command
15 9	* * *	bulovic	/home/bulovic/NetBeansProjects/Lake-Water-Level/scrape.pl
